# scripts/polyclawd-monitor.service
[Unit]
Description=Polyclawd Signal Monitor
After=network.target polyclawd-api.service

[Service]
Type=simple
User=linuxuser
WorkingDirectory=/var/www/virtuosocrypto.com/polyclawd
ExecStart=/usr/bin/python3 /var/www/virtuosocrypto.com/polyclawd/scripts/monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
